# etc/cron.d/vulners
0 3 * * * root /usr/local/bin/cvescan --priority all --json > /var/log/vulners_report; /usr/bin/python3 /usr/local/bin/cvescan_handler -o file -rp /var/log/vulners_report
